# raspios/base/etc/rc.local
#!/bin/sh
#
# rc.local

# run setup scripts
if test -d /etc/boot.d
then
    run-parts --lsbsysinit /etc/boot.d
fi
if test -d /usr/local/etc/boot.d
then
    run-parts --lsbsysinit /usr/local/etc/boot.d
fi

exit 0
# End
